# lib/systemd/system/acpid.service
[Unit]
Description=ACPI event daemon
Requires=acpid.socket

[Service]
StandardInput=socket
ExecStart=/usr/sbin/acpid
