# lib/systemd/system/run-lock.mount
[Unit]
Description=Lock Directory
Before=local-fs.target

[Mount]
What=tmpfs
Where=/run/lock
Type=tmpfs
Options=nodev,noexec,nosuid,size=5242880
