# lib/systemd/system/startx.service
[Unit]
Description=Manages X

[Service]
ExecStart=/etc/webc/startx.sh

[Install]
WantedBy=graphical.target
