# lib/systemd/system/startx.service
[Unit]
Description=Manages X and replaces getty

[Service]
ExecStart=/etc/webc/startx.sh
Restart=always

[Install]
WantedBy=multi-user.target
