# var/lib/update-rc.d/rc.local
update-rc.d rc.local start 99 2 3 4 5 .
