# lib/udev/rules.d/80-alsa.rules
KERNEL=="controlC[0-7]", ACTION=="add", RUN+="/lib/udev/alsa-utils"
