# lib/udev/rules.d/80-networking.rules
SUBSYSTEM=="net",			RUN+="net.agent"

SUBSYSTEM=="atm",			RUN+="dsl-modem.agent"

