# lib/udev/rules.d/90-alsa-restore.rules
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="controlC*", KERNELS=="card*", \
        TEST=="/usr/sbin/alsactl", RUN+="/usr/sbin/alsactl restore $attr{number}"
