# var/lib/dpkg/info/acpi-support-base.postinst
#!/bin/sh

set -e

case "$1" in
        configure)		
                invoke-rc.d acpid restart >/dev/null || true
                ;;
        *)
                ;;
esac


