# var/lib/dpkg/info/acpid.list
/.
/etc
/etc/acpi
/etc/acpi/events
/etc/default
/etc/default/acpid
/etc/init.d
/etc/init.d/acpid
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/acpid.path
/lib/systemd/system/acpid.service
/lib/systemd/system/acpid.socket
/usr
/usr/bin
/usr/bin/acpi_listen
/usr/sbin
/usr/sbin/acpid
/usr/share
/usr/share/doc
/usr/share/doc/acpid
/usr/share/doc/acpid/README
/usr/share/doc/acpid/README.Debian
/usr/share/doc/acpid/changelog.Debian.gz
/usr/share/doc/acpid/changelog.Debian.i386.gz
/usr/share/doc/acpid/changelog.gz
/usr/share/doc/acpid/copyright
/usr/share/doc/acpid/examples
/usr/share/doc/acpid/examples/ac
/usr/share/doc/acpid/examples/ac.sh
/usr/share/doc/acpid/examples/default
/usr/share/doc/acpid/examples/default.sh
/usr/share/doc/acpid/examples/powerbtn
/usr/share/doc/acpid/examples/powerbtn.sh
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/acpi_listen.8.gz
/usr/share/man/man8/acpid.8.gz
