# var/lib/dpkg/info/alsa-utils.list
/.
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-alsa-restore.rules
/etc
/etc/init.d
/etc/init.d/alsa-utils
/var
/var/lib
/var/lib/alsa
/usr
/usr/bin
/usr/bin/aplay
/usr/bin/amixer
/usr/bin/arecordmidi
/usr/bin/aseqnet
/usr/bin/alsamixer
/usr/bin/alsaucm
/usr/bin/aconnect
/usr/bin/speaker-test
/usr/bin/iecset
/usr/bin/amidi
/usr/bin/alsaloop
/usr/bin/aplaymidi
/usr/bin/aseqdump
/usr/sbin
/usr/sbin/alsactl
/usr/share
/usr/share/doc
/usr/share/doc/alsa-utils
/usr/share/doc/alsa-utils/changelog.gz
/usr/share/doc/alsa-utils/changelog.Debian.gz
/usr/share/doc/alsa-utils/copyright
/usr/share/doc/alsa-utils/README.Debian
/usr/share/doc/alsa-utils/NEWS.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/alsa-utils
/usr/share/locale/ja/LC_MESSAGES/alsa-utils.mo
/usr/share/locale/fr/LC_MESSAGES/alsa-utils.mo
/usr/share/locale/de/LC_MESSAGES/alsa-utils.mo
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/alsactl_init.7.gz
/usr/share/man/man1
/usr/share/man/man1/aplay.1.gz
/usr/share/man/man1/aseqdump.1.gz
/usr/share/man/man1/speaker-test.1.gz
/usr/share/man/man1/aconnect.1.gz
/usr/share/man/man1/aplaymidi.1.gz
/usr/share/man/man1/aseqnet.1.gz
/usr/share/man/man1/arecordmidi.1.gz
/usr/share/man/man1/alsactl.1.gz
/usr/share/man/man1/amidi.1.gz
/usr/share/man/man1/alsaloop.1.gz
/usr/share/man/man1/iecset.1.gz
/usr/share/man/man1/amixer.1.gz
/usr/share/man/man1/alsamixer.1.gz
/usr/share/sounds
/usr/share/sounds/alsa
/usr/share/sounds/alsa/Side_Left.wav
/usr/share/sounds/alsa/Noise.wav
/usr/share/sounds/alsa/Front_Left.wav
/usr/share/sounds/alsa/Front_Center.wav
/usr/share/sounds/alsa/Rear_Left.wav
/usr/share/sounds/alsa/Rear_Center.wav
/usr/share/sounds/alsa/Front_Right.wav
/usr/share/sounds/alsa/Rear_Right.wav
/usr/share/sounds/alsa/Side_Right.wav
/usr/share/alsa
/usr/share/alsa/utils.sh
/usr/share/alsa/speaker-test
/usr/share/alsa/speaker-test/sample_map.csv
/usr/share/alsa/init
/usr/share/alsa/init/test
/usr/share/alsa/init/hda
/usr/share/alsa/init/help
/usr/share/alsa/init/00main
/usr/share/alsa/init/info
/usr/share/alsa/init/default
/usr/bin/arecord
/usr/share/man/man1/arecord.1.gz
