# var/lib/dpkg/info/aptitude.conffiles
/etc/cron.daily/aptitude
/etc/logrotate.d/aptitude
