# var/lib/dpkg/info/aptitude.list
/etc
/etc/cron.daily
/etc/cron.daily/aptitude
/etc/logrotate.d
/etc/logrotate.d/aptitude
/var/lib/aptitude
