# var/lib/dpkg/info/consolekit.list
/.
/etc
/etc/ConsoleKit
/etc/ConsoleKit/run-seat.d
/etc/ConsoleKit/run-session.d
/etc/ConsoleKit/seats.d
/etc/ConsoleKit/seats.d/00-primary.seat
/etc/X11
/etc/X11/Xsession.d
/etc/X11/Xsession.d/90consolekit
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/ConsoleKit.conf
/etc/logrotate.d
/etc/logrotate.d/consolekit
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/basic.target.wants
/lib/systemd/system/console-kit-daemon.service
/lib/systemd/system/console-kit-log-system-restart.service
/lib/systemd/system/console-kit-log-system-start.service
/lib/systemd/system/console-kit-log-system-stop.service
/lib/systemd/system/halt.target.wants
/lib/systemd/system/poweroff.target.wants
/lib/systemd/system/reboot.target.wants
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/70-udev-acl.rules
/lib/udev/udev-acl
/usr
/usr/bin
/usr/bin/ck-history
/usr/bin/ck-launch-session
/usr/bin/ck-list-sessions
/usr/lib
/usr/lib/ConsoleKit
/usr/lib/ConsoleKit/ck-collect-session-info
/usr/lib/ConsoleKit/ck-get-x11-display-device
/usr/lib/ConsoleKit/ck-get-x11-server-pid
/usr/lib/ConsoleKit/run-seat.d
/usr/lib/ConsoleKit/run-session.d
/usr/lib/ConsoleKit/run-session.d/pam-foreground-compat.ck
/usr/lib/ConsoleKit/scripts
/usr/lib/ConsoleKit/scripts/ck-system-restart
/usr/lib/ConsoleKit/scripts/ck-system-stop
/usr/sbin
/usr/sbin/ck-log-system-restart
/usr/sbin/ck-log-system-start
/usr/sbin/ck-log-system-stop
/usr/sbin/console-kit-daemon
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.ConsoleKit.service
/usr/share/doc
/usr/share/doc/consolekit
/usr/share/doc/consolekit/AUTHORS
/usr/share/doc/consolekit/NEWS.gz
/usr/share/doc/consolekit/README
/usr/share/doc/consolekit/README.Debian
/usr/share/doc/consolekit/TODO
/usr/share/doc/consolekit/changelog.Debian.gz
/usr/share/doc/consolekit/changelog.gz
/usr/share/doc/consolekit/copyright
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.consolekit.policy
/var
/var/log
/var/log/ConsoleKit
/lib/systemd/system/basic.target.wants/console-kit-log-system-start.service
/lib/systemd/system/halt.target.wants/console-kit-log-system-stop.service
/lib/systemd/system/poweroff.target.wants/console-kit-log-system-stop.service
/lib/systemd/system/reboot.target.wants/console-kit-log-system-restart.service
/usr/lib/ConsoleKit/run-seat.d/udev-acl.ck
