# var/lib/dpkg/info/cups.conffiles
/etc/logrotate.d/cups
/etc/pam.d/cups
/etc/init.d/cups
/etc/cups/cupsd.conf
/etc/cups/snmp.conf
/etc/cups/cupsd.conf.default
/etc/default/cups
