# var/lib/dpkg/info/dbus.conffiles
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
/etc/init.d/dbus
/etc/default/dbus
