# var/lib/dpkg/info/google-talkplugin.list
/.
/opt
/opt/google
/opt/google/talkplugin
/opt/google/talkplugin/libgoogletalkremoting.so
/opt/google/talkplugin/GoogleTalkPlugin
/opt/google/talkplugin/attributions.txt
/opt/google/talkplugin/lib
/opt/google/talkplugin/lib/libCgGL.so
/opt/google/talkplugin/lib/libCg.so
/opt/google/talkplugin/remoting24x24.png
/opt/google/talkplugin/libnpo1d.so
/opt/google/talkplugin/cron
/opt/google/talkplugin/cron/google-talkplugin
/opt/google/talkplugin/windowpicker.glade
/opt/google/talkplugin/libnpgoogletalk.so
/opt/google/talkplugin/libnpgtpo3dautoplugin.so
/opt/google/talkplugin/data
/opt/google/talkplugin/data/SFTprec_120.emd
/opt/google/talkplugin/data/LMprec_508.emd
/opt/google/talkplugin/locale
/opt/google/talkplugin/locale/fil
/opt/google/talkplugin/locale/fil/LC_MESSAGES
/opt/google/talkplugin/locale/fil/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/nl
/opt/google/talkplugin/locale/nl/LC_MESSAGES
/opt/google/talkplugin/locale/nl/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/kn
/opt/google/talkplugin/locale/kn/LC_MESSAGES
/opt/google/talkplugin/locale/kn/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/da
/opt/google/talkplugin/locale/da/LC_MESSAGES
/opt/google/talkplugin/locale/da/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/es-419
/opt/google/talkplugin/locale/es-419/LC_MESSAGES
/opt/google/talkplugin/locale/es-419/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/sk
/opt/google/talkplugin/locale/sk/LC_MESSAGES
/opt/google/talkplugin/locale/sk/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/fr-CA
/opt/google/talkplugin/locale/fr-CA/LC_MESSAGES
/opt/google/talkplugin/locale/fr-CA/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/vi
/opt/google/talkplugin/locale/vi/LC_MESSAGES
/opt/google/talkplugin/locale/vi/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/pl
/opt/google/talkplugin/locale/pl/LC_MESSAGES
/opt/google/talkplugin/locale/pl/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ja
/opt/google/talkplugin/locale/ja/LC_MESSAGES
/opt/google/talkplugin/locale/ja/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/sr
/opt/google/talkplugin/locale/sr/LC_MESSAGES
/opt/google/talkplugin/locale/sr/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/zh-HK
/opt/google/talkplugin/locale/zh-HK/LC_MESSAGES
/opt/google/talkplugin/locale/zh-HK/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ro
/opt/google/talkplugin/locale/ro/LC_MESSAGES
/opt/google/talkplugin/locale/ro/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/bg
/opt/google/talkplugin/locale/bg/LC_MESSAGES
/opt/google/talkplugin/locale/bg/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/zh-CN
/opt/google/talkplugin/locale/zh-CN/LC_MESSAGES
/opt/google/talkplugin/locale/zh-CN/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/sw
/opt/google/talkplugin/locale/sw/LC_MESSAGES
/opt/google/talkplugin/locale/sw/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/no
/opt/google/talkplugin/locale/no/LC_MESSAGES
/opt/google/talkplugin/locale/no/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/hr
/opt/google/talkplugin/locale/hr/LC_MESSAGES
/opt/google/talkplugin/locale/hr/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ru
/opt/google/talkplugin/locale/ru/LC_MESSAGES
/opt/google/talkplugin/locale/ru/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/pt-PT
/opt/google/talkplugin/locale/pt-PT/LC_MESSAGES
/opt/google/talkplugin/locale/pt-PT/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/id
/opt/google/talkplugin/locale/id/LC_MESSAGES
/opt/google/talkplugin/locale/id/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ms
/opt/google/talkplugin/locale/ms/LC_MESSAGES
/opt/google/talkplugin/locale/ms/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/hu
/opt/google/talkplugin/locale/hu/LC_MESSAGES
/opt/google/talkplugin/locale/hu/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/zu
/opt/google/talkplugin/locale/zu/LC_MESSAGES
/opt/google/talkplugin/locale/zu/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/sv
/opt/google/talkplugin/locale/sv/LC_MESSAGES
/opt/google/talkplugin/locale/sv/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/cs
/opt/google/talkplugin/locale/cs/LC_MESSAGES
/opt/google/talkplugin/locale/cs/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/gu
/opt/google/talkplugin/locale/gu/LC_MESSAGES
/opt/google/talkplugin/locale/gu/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/is
/opt/google/talkplugin/locale/is/LC_MESSAGES
/opt/google/talkplugin/locale/is/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/fa
/opt/google/talkplugin/locale/fa/LC_MESSAGES
/opt/google/talkplugin/locale/fa/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/fi
/opt/google/talkplugin/locale/fi/LC_MESSAGES
/opt/google/talkplugin/locale/fi/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/or
/opt/google/talkplugin/locale/or/LC_MESSAGES
/opt/google/talkplugin/locale/or/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/af
/opt/google/talkplugin/locale/af/LC_MESSAGES
/opt/google/talkplugin/locale/af/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/el
/opt/google/talkplugin/locale/el/LC_MESSAGES
/opt/google/talkplugin/locale/el/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/de
/opt/google/talkplugin/locale/de/LC_MESSAGES
/opt/google/talkplugin/locale/de/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/te
/opt/google/talkplugin/locale/te/LC_MESSAGES
/opt/google/talkplugin/locale/te/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/it
/opt/google/talkplugin/locale/it/LC_MESSAGES
/opt/google/talkplugin/locale/it/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/lv
/opt/google/talkplugin/locale/lv/LC_MESSAGES
/opt/google/talkplugin/locale/lv/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/tl
/opt/google/talkplugin/locale/tl/LC_MESSAGES
/opt/google/talkplugin/locale/tl/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/pt-BR
/opt/google/talkplugin/locale/pt-BR/LC_MESSAGES
/opt/google/talkplugin/locale/pt-BR/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ko
/opt/google/talkplugin/locale/ko/LC_MESSAGES
/opt/google/talkplugin/locale/ko/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/en-GB
/opt/google/talkplugin/locale/en-GB/LC_MESSAGES
/opt/google/talkplugin/locale/en-GB/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ca
/opt/google/talkplugin/locale/ca/LC_MESSAGES
/opt/google/talkplugin/locale/ca/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ur
/opt/google/talkplugin/locale/ur/LC_MESSAGES
/opt/google/talkplugin/locale/ur/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ml
/opt/google/talkplugin/locale/ml/LC_MESSAGES
/opt/google/talkplugin/locale/ml/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/lt
/opt/google/talkplugin/locale/lt/LC_MESSAGES
/opt/google/talkplugin/locale/lt/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/zh-TW
/opt/google/talkplugin/locale/zh-TW/LC_MESSAGES
/opt/google/talkplugin/locale/zh-TW/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/et
/opt/google/talkplugin/locale/et/LC_MESSAGES
/opt/google/talkplugin/locale/et/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/en
/opt/google/talkplugin/locale/en/LC_MESSAGES
/opt/google/talkplugin/locale/en/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/eu
/opt/google/talkplugin/locale/eu/LC_MESSAGES
/opt/google/talkplugin/locale/eu/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/es
/opt/google/talkplugin/locale/es/LC_MESSAGES
/opt/google/talkplugin/locale/es/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/am
/opt/google/talkplugin/locale/am/LC_MESSAGES
/opt/google/talkplugin/locale/am/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ar
/opt/google/talkplugin/locale/ar/LC_MESSAGES
/opt/google/talkplugin/locale/ar/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/gl
/opt/google/talkplugin/locale/gl/LC_MESSAGES
/opt/google/talkplugin/locale/gl/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/ta
/opt/google/talkplugin/locale/ta/LC_MESSAGES
/opt/google/talkplugin/locale/ta/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/tr
/opt/google/talkplugin/locale/tr/LC_MESSAGES
/opt/google/talkplugin/locale/tr/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/iw
/opt/google/talkplugin/locale/iw/LC_MESSAGES
/opt/google/talkplugin/locale/iw/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/bn
/opt/google/talkplugin/locale/bn/LC_MESSAGES
/opt/google/talkplugin/locale/bn/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/fr
/opt/google/talkplugin/locale/fr/LC_MESSAGES
/opt/google/talkplugin/locale/fr/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/uk
/opt/google/talkplugin/locale/uk/LC_MESSAGES
/opt/google/talkplugin/locale/uk/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/th
/opt/google/talkplugin/locale/th/LC_MESSAGES
/opt/google/talkplugin/locale/th/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/mr
/opt/google/talkplugin/locale/mr/LC_MESSAGES
/opt/google/talkplugin/locale/mr/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/sl
/opt/google/talkplugin/locale/sl/LC_MESSAGES
/opt/google/talkplugin/locale/sl/LC_MESSAGES/windowpicker.mo
/opt/google/talkplugin/locale/hi
/opt/google/talkplugin/locale/hi/LC_MESSAGES
/opt/google/talkplugin/locale/hi/LC_MESSAGES/windowpicker.mo
/etc
/etc/cron.daily
/usr
/usr/share
/usr/share/doc
/usr/share/doc/google-talkplugin
/usr/share/doc/google-talkplugin/changelog.Debian.gz
/usr/lib
/usr/lib/xulrunner
/usr/lib/xulrunner/plugins
/usr/lib/firefox
/usr/lib/firefox/plugins
/usr/lib/iceape
/usr/lib/iceape/plugins
/usr/lib/iceweasel
/usr/lib/iceweasel/plugins
/usr/lib/midbrowser
/usr/lib/midbrowser/plugins
/usr/lib/xulrunner-addons
/usr/lib/xulrunner-addons/plugins
/usr/lib/mozilla
/usr/lib/mozilla/plugins
/etc/cron.daily/google-talkplugin
/usr/lib/xulrunner/plugins/libnpo1d.so
/usr/lib/xulrunner/plugins/libnpgoogletalk.so
/usr/lib/xulrunner/plugins/libnpgtpo3dautoplugin.so
/usr/lib/firefox/plugins/libnpo1d.so
/usr/lib/firefox/plugins/libnpgoogletalk.so
/usr/lib/firefox/plugins/libnpgtpo3dautoplugin.so
/usr/lib/iceape/plugins/libnpo1d.so
/usr/lib/iceape/plugins/libnpgoogletalk.so
/usr/lib/iceape/plugins/libnpgtpo3dautoplugin.so
/usr/lib/iceweasel/plugins/libnpo1d.so
/usr/lib/iceweasel/plugins/libnpgoogletalk.so
/usr/lib/iceweasel/plugins/libnpgtpo3dautoplugin.so
/usr/lib/midbrowser/plugins/libnpo1d.so
/usr/lib/midbrowser/plugins/libnpgoogletalk.so
/usr/lib/midbrowser/plugins/libnpgtpo3dautoplugin.so
/usr/lib/xulrunner-addons/plugins/libnpo1d.so
/usr/lib/xulrunner-addons/plugins/libnpgoogletalk.so
/usr/lib/xulrunner-addons/plugins/libnpgtpo3dautoplugin.so
/usr/lib/mozilla/plugins/libnpo1d.so
/usr/lib/mozilla/plugins/libnpgoogletalk.so
/usr/lib/mozilla/plugins/libnpgtpo3dautoplugin.so
