# var/lib/dpkg/info/ifupdown.conffiles
/etc/default/networking
/etc/init.d/networking
/etc/init/network-interface-container.conf
/etc/init/network-interface-security.conf
/etc/init/network-interface.conf
/etc/init/networking.conf
/etc/network/if-down.d/upstart
/etc/network/if-up.d/upstart
