# var/lib/dpkg/info/ifupdown.list
/.
/sbin
/sbin/ifdown
/sbin/ifup
/usr
/usr/share
/usr/share/ifupdown
/usr/share/ifupdown/upgrade-from-0.5.x.pl
/usr/share/man/man5/interfaces.5.gz
/usr/share/man/man8/ifup.8.gz
/usr/share/doc/ifupdown
/usr/share/doc/ifupdown/changelog.Debian.gz
/usr/share/doc/ifupdown/contrib
/usr/share/doc/ifupdown/contrib/ensureifup
/usr/share/doc/ifupdown/contrib/ifstate-check
/usr/share/doc/ifupdown/contrib/ifstate
/usr/share/doc/ifupdown/TODO
/usr/share/doc/ifupdown/examples
/usr/share/doc/ifupdown/examples/check-mac-address.sh
/usr/share/doc/ifupdown/examples/ping-places.sh
/usr/share/doc/ifupdown/examples/bridge
/usr/share/doc/ifupdown/examples/generate-interfaces.pl.gz
/usr/share/doc/ifupdown/examples/pcmcia-compat.sh
/usr/share/doc/ifupdown/examples/get-mac-address.sh
/usr/share/doc/ifupdown/examples/network-interfaces.gz
/usr/share/doc/ifupdown/copyright
/usr/share/doc/ifupdown/changelog.gz
/etc
/etc/network
/etc/network/if-pre-up.d
/etc/network/if-down.d
/etc/network/if-up.d
/etc/network/if-post-down.d
/etc/default
/etc/default/ifupdown
/etc/init.d
/etc/init.d/ifupdown
/etc/init.d/ifupdown-clean
/usr/share/man/man8/ifdown.8.gz
