# var/lib/dpkg/info/kbd.conffiles
/etc/kbd/remap
/etc/kbd/config
/etc/init.d/kbd
