# var/lib/dpkg/info/keyboard-configuration.conffiles
/etc/init.d/console-setup
/etc/init.d/keyboard-setup
