# var/lib/dpkg/info/keyboard-configuration.list
/.
/usr
/usr/share
/usr/share/console-setup
/usr/share/console-setup/keyboard
/usr/share/doc
/usr/share/doc/keyboard-configuration
/usr/share/doc/keyboard-configuration/README.Debian
/usr/share/doc/keyboard-configuration/copyright
/usr/share/doc/keyboard-configuration/examples
/usr/share/doc/keyboard-configuration/changelog.gz
/etc
/etc/init.d
/etc/init.d/console-setup
/etc/init.d/keyboard-setup
/usr/share/doc/keyboard-configuration/examples/keyboard
