# var/lib/dpkg/info/libdevmapper1.02.1.preinst
#!/bin/sh

if [ "$1" = "upgrade" ] && dpkg --compare-versions "$2" "<<" "2:1.02.24"; then
	update-rc.d -f libdevmapper1.02.1 remove >/dev/null || exit $?
	rm /etc/init.d/libdevmapper1.02.1
fi
if [ "$1" = "upgrade" ] && dpkg --compare-versions "$2" "<=" "2:1.02.48-4"; then
	update-rc.d -f libdevmapper1.02 remove >/dev/null
	rm -f /etc/init.d/libdevmapper1.02
fi



exit 0
