# var/lib/dpkg/info/libgl1-mesa-dri:i386.list
/.
/usr
/usr/lib
/usr/lib/i386-linux-gnu
/usr/lib/i386-linux-gnu/dri
/usr/lib/i386-linux-gnu/dri/i915_dri.so
/usr/lib/i386-linux-gnu/dri/i965_dri.so
/usr/lib/i386-linux-gnu/dri/kms_swrast_dri.so
/usr/lib/i386-linux-gnu/dri/nouveau_dri.so
/usr/lib/i386-linux-gnu/dri/nouveau_vieux_dri.so
/usr/lib/i386-linux-gnu/dri/r200_dri.so
/usr/lib/i386-linux-gnu/dri/r300_dri.so
/usr/lib/i386-linux-gnu/dri/r600_dri.so
/usr/lib/i386-linux-gnu/dri/radeon_dri.so
/usr/lib/i386-linux-gnu/dri/radeonsi_dri.so
/usr/lib/i386-linux-gnu/dri/swrast_dri.so
/usr/lib/i386-linux-gnu/dri/virtio_gpu_dri.so
/usr/lib/i386-linux-gnu/dri/vmwgfx_dri.so
/usr/share
/usr/share/bug
/usr/share/bug/libgl1-mesa-dri
/usr/share/bug/libgl1-mesa-dri/control
/usr/share/bug/libgl1-mesa-dri/script
/usr/share/doc
/usr/share/doc/libgl1-mesa-dri
/usr/share/doc/libgl1-mesa-dri/changelog.Debian.gz
/usr/share/doc/libgl1-mesa-dri/copyright
/usr/share/drirc.d
/usr/share/drirc.d/00-mesa-defaults.conf
