# var/lib/dpkg/info/live-config-sysvinit.list
/.
/usr
/usr/share
/usr/share/doc/live-config-sysvinit
/usr/share/doc/live-config-sysvinit/changelog.Debian.gz
/usr/share/doc/live-config-sysvinit/copyright
/lib
/lib/live
/lib/live/config
/lib/live/config/017-sysv-rc
/lib/live/config/016-sysvinit
/etc
/etc/init.d
/etc/init.d/live-config
