# var/lib/dpkg/info/live-config-systemd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/live-config.service
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-config-systemd
/usr/share/doc/live-config-systemd/copyright
/usr/share/doc/live-config-systemd/changelog.Debian.gz
/etc
/etc/systemd
/etc/systemd/system
/etc/systemd/system/basic.target.wants
/etc/systemd/system/basic.target.wants/live-config.service
