# var/lib/dpkg/info/live-config-systemd.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-config-systemd
/usr/share/doc/live-config-systemd/changelog.Debian.gz
/usr/share/doc/live-config-systemd/copyright
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/live-config.service
