# var/lib/dpkg/info/live-boot.list
/.
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/live-boot
/usr/share/doc
/usr/share/doc/live-boot
/usr/share/doc/live-boot/old
/usr/share/doc/live-boot/old/parameters.txt
/usr/share/doc/live-boot/old/ChangeLog.casper.gz
/usr/share/doc/live-boot/old/CREDITS
/usr/share/doc/live-boot/old/ChangeLog.gz
/usr/share/doc/live-boot/examples
/usr/share/doc/live-boot/examples/live-snapshot.exclude_list
/usr/share/doc/live-boot/examples/live-snapshot.list
/usr/share/doc/live-boot/examples/live-persistence.binds
/usr/share/doc/live-boot/copyright
/usr/share/doc/live-boot/changelog.Debian.gz
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/live-boot.7.gz
/usr/share/man/de
/usr/share/man/de/man7
/usr/share/man/de/man7/live-boot.7.gz
/usr/share/man/de/man1
/usr/share/man/de/man1/live-snapshot.1.gz
/usr/share/man/man1
/usr/share/man/man1/live-snapshot.1.gz
/usr/share/bug
/usr/share/bug/live-boot
/usr/share/bug/live-boot/script
/usr/share/bug/live-boot/presubj
/usr/share/live-boot
/usr/share/live-boot/live-preseed
/usr/share/live-boot/languagelist
/usr/share/live-boot/live-reconfigure
/bin
/bin/live-system
/etc
/etc/init.d
/etc/init.d/live-boot
/sbin
/sbin/live-swapfile
/sbin/live-new-uuid
/sbin/live-toram
/sbin/live-snapshot
