# var/lib/dpkg/info/live-config-sysvinit.conffiles
/etc/init.d/live-config
/etc/init.d/live
