# var/lib/dpkg/info/live-config-sysvinit.list
/lib
/lib/live
/lib/live/config
/etc
/etc/init.d
/etc/init.d/live-config
/etc/init.d/live
