# var/lib/dpkg/info/live-boot.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-boot
/usr/share/doc/live-boot/changelog.Debian.gz
/usr/share/doc/live-boot/copyright
/usr/share/bug
/usr/share/bug/live-boot
/usr/share/bug/live-boot/presubj
/usr/share/bug/live-boot/script
/lib
/lib/live
/lib/live/boot.sh
/lib/live/boot
/lib/live/boot/swapon.sh
/lib/live/boot/mount-cifs.sh
/lib/live/boot/netboot.sh
/lib/live/boot/mount-nfs.sh
/lib/live/boot/toram-todisk.sh
/lib/live/boot/read-only.sh
/lib/live/boot/networking.sh
/lib/live/boot/overlay.sh
/lib/live/boot/fstab.sh
/lib/live/boot/mount-http.sh
/lib/live/boot/select-eth-device.sh
/lib/live/boot/verify-checksums.sh
/lib/live/boot/initramfs-tools.sh
/lib/live/boot/cmdline.sh
/lib/live/boot/mount-iscsi.sh
/lib/live/boot/netbase.sh
/lib/live/boot/misc-helpers.sh
/lib/live/boot/mountroot.sh
/sbin
/sbin/live-swapfile
/sbin/live-new-uuid
/etc/init.d/live-boot
