# var/lib/dpkg/info/live-config-systemd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/live-config.service
/lib/systemd/system-generators
/lib/systemd/system-generators/live-config-getty-generator
/usr
/usr/share
/usr/share/doc
/usr/share/doc/live-config-systemd
/usr/share/doc/live-config-systemd/changelog.gz
/usr/share/doc/live-config-systemd/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/live-config-systemd
