# var/lib/dpkg/info/logrotate.list
/.
/usr
/usr/sbin
/usr/sbin/logrotate
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/logrotate.8.gz
/usr/share/bug
/usr/share/bug/logrotate
/usr/share/bug/logrotate/script
/usr/share/doc
/usr/share/doc/logrotate
/usr/share/doc/logrotate/README.Debian
/usr/share/doc/logrotate/copyright
/usr/share/doc/logrotate/changelog.Debian.gz
/usr/share/doc/logrotate/changelog.gz
/var
/var/lib
/var/lib/logrotate
/etc
/etc/logrotate.d
/etc/logrotate.conf
/etc/cron.daily
/etc/cron.daily/logrotate
