# var/lib/dpkg/info/man-db.conffiles
/etc/cron.daily/man-db
/etc/cron.weekly/man-db
/etc/manpath.config
