# var/lib/dpkg/info/man-db.conffiles
/etc/apparmor.d/usr.bin.man
/etc/cron.daily/man-db
/etc/cron.weekly/man-db
/etc/manpath.config
