# var/lib/dpkg/info/module-init-tools.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/module-init-tools
/usr/share/doc/module-init-tools/changelog.gz
/usr/share/doc/module-init-tools/NEWS.gz
/usr/share/doc/module-init-tools/changelog.Debian.gz
/usr/share/doc/module-init-tools/copyright
/usr/share/man
/usr/share/man/fr
/usr/share/man/fr/man5
/usr/share/man/fr/man5/modules.5.gz
/usr/share/man/man8
/usr/share/man/man8/depmod.8.gz
/usr/share/man/man8/modinfo.8.gz
/usr/share/man/man8/lsmod.8.gz
/usr/share/man/man8/rmmod.8.gz
/usr/share/man/man8/insmod.8.gz
/usr/share/man/man8/modprobe.8.gz
/usr/share/man/man5
/usr/share/man/man5/modules.dep.5.gz
/usr/share/man/man5/depmod.d.5.gz
/usr/share/man/man5/modules.dep.bin.5.gz
/usr/share/man/man5/modprobe.conf.5.gz
/usr/share/man/man5/depmod.conf.5.gz
/usr/share/man/man5/modules.5.gz
/bin
/bin/lsmod
/lib
/lib/modules
/sbin
/sbin/modprobe
/sbin/rmmod
/sbin/modinfo
/sbin/insmod
/sbin/depmod
/etc
/etc/init.d
/etc/init.d/module-init-tools
/etc/modprobe.d
/etc/modprobe.d/aliases.conf
/sbin/lsmod
