# var/lib/dpkg/info/mozplugger.conffiles
/etc/mozpluggerrc.d/20-header.conf
/etc/mozpluggerrc.d/40-global_defs.conf
/etc/mozpluggerrc.d/60-video.conf
/etc/mozpluggerrc.d/61-audio.conf
/etc/mozpluggerrc.d/62-documents.conf
