# var/lib/dpkg/info/mozplugger.list
/.
/etc
/etc/mozpluggerrc.d
/etc/mozpluggerrc.d/20-header.conf
/etc/mozpluggerrc.d/40-global_defs.conf
/etc/mozpluggerrc.d/60-video.conf
/etc/mozpluggerrc.d/61-audio.conf
/etc/mozpluggerrc.d/62-documents.conf
/usr
/usr/bin
/usr/bin/mozplugger-controller
/usr/bin/mozplugger-helper
/usr/bin/mozplugger-linker
/usr/lib
/usr/lib/mozilla
/usr/lib/mozilla/plugins
/usr/lib/mozilla/plugins/mozplugger.so
/usr/sbin
/usr/sbin/update-mozpluggerrc
/usr/share
/usr/share/doc
/usr/share/doc/mozplugger
/usr/share/doc/mozplugger/NEWS.Debian.gz
/usr/share/doc/mozplugger/README
/usr/share/doc/mozplugger/README.Debian
/usr/share/doc/mozplugger/changelog.Debian.gz
/usr/share/doc/mozplugger/changelog.Debian.i386.gz
/usr/share/doc/mozplugger/changelog.gz
/usr/share/doc/mozplugger/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mozplugger-controller.1.gz
/usr/share/man/man1/mozplugger-helper.1.gz
/usr/share/man/man7
/usr/share/man/man7/mozplugger.7.gz
/usr/share/man/man8
/usr/share/man/man8/update-mozpluggerrc.8.gz
/usr/share/man/man1/mozplugger-linker.1.gz
