# var/lib/dpkg/info/netbase.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/netbase
/usr/share/doc/netbase/changelog.gz
/usr/share/doc/netbase/README.Debian
/usr/share/doc/netbase/copyright
/etc
/etc/network
/etc/services
/etc/protocols
/etc/rpc
/etc/init.d
/etc/init.d/networking
