# var/lib/dpkg/info/nvidia-kernel-common.conffiles
/etc/init.d/nvidia-kernel
/etc/modprobe.d/nvidia-kernel-common.conf
/etc/default/nvidia-kernel
