# var/lib/dpkg/info/nvidia-kernel-common.list
/etc
/etc/init.d
/etc/init.d/nvidia-kernel
/etc/modprobe.d
/etc/modprobe.d/nvidia-kernel-common.conf
/etc/default
/etc/default/nvidia-kernel
