# var/lib/dpkg/info/nvidia-kernel-common.list
/.
/usr
/usr/share
/usr/share/doc
/usr/share/doc/nvidia-kernel-common
/usr/share/doc/nvidia-kernel-common/changelog.Debian.gz
/usr/share/doc/nvidia-kernel-common/README.Debian
/usr/share/doc/nvidia-kernel-common/copyright
/usr/lib
/usr/lib/ConsoleKit
/usr/lib/ConsoleKit/run-seat.d
/usr/lib/ConsoleKit/run-seat.d/nvidia_helper.ck
/lib
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/60-nvidia-kernel-common.rules
/lib/udev/nvidia_helper
/etc
/etc/init.d
/etc/init.d/nvidia-kernel
/etc/modprobe.d
/etc/modprobe.d/nvidia-kernel-common.conf
/etc/default
/etc/default/nvidia-kernel
