# var/lib/dpkg/info/plymouth.list
/.
/lib
/lib/libply.so.2.0.0
/lib/libply-splash-core.so.2.0.0
/bin
/bin/plymouth
/usr
/usr/lib
/usr/lib/libply-boot-client.so.2.0.0
/usr/lib/plymouth
/usr/lib/plymouth/plymouth
/usr/lib/plymouth/plymouth/plymouth-populate-initrd
/usr/lib/plymouth/plymouth/plymouth-update-initrd
/usr/lib/plymouth/plymouth/plymouth-generate-initrd
/usr/lib/plymouth/text.so
/usr/lib/plymouth/details.so
/usr/share
/usr/share/man/man8
/usr/share/man/man8/plymouth.8.gz
/usr/share/plymouth
/usr/share/plymouth/themes
/usr/share/plymouth/themes/text
/usr/share/plymouth/themes/text/text.plymouth
/usr/share/plymouth/themes/details
/usr/share/plymouth/themes/details/details.plymouth
/usr/share/plymouth/plymouthd.defaults
/usr/share/plymouth/debian-logo.png
/usr/share/doc/plymouth
/usr/share/doc/plymouth/README.Debian
/usr/share/doc/plymouth/changelog.Debian.gz
/usr/share/doc/plymouth/TODO
/usr/share/doc/plymouth/README
/usr/share/doc/plymouth/changelog.gz
/usr/share/doc/plymouth/copyright
/usr/share/lintian/overrides/plymouth
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/plymouth
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/sbin
/usr/sbin/plymouth-set-default-theme
/sbin
/sbin/plymouthd
/etc
/etc/plymouth
/etc/plymouth/plymouthd.conf
/etc/init.d
/etc/init.d/plymouth
/lib/libply-splash-core.so.2
/lib/libply.so.2
/usr/lib/libply-boot-client.so.2
/usr/share/man/man8/plymouth-set-default-theme.8.gz
