# var/lib/dpkg/info/plymouth.list
/.
/bin
/bin/plymouth
/etc
/etc/init.d
/etc/init.d/plymouth
/etc/init.d/plymouth-log
/etc/plymouth
/etc/plymouth/plymouthd.conf
/lib
/lib/lsb
/lib/lsb/init-functions.d
/lib/lsb/init-functions.d/99-plymouth
/lib/systemd
/lib/systemd/system
/lib/systemd/system/halt.target.wants
/lib/systemd/system/initrd-switch-root.target.wants
/lib/systemd/system/kexec.target.wants
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/plymouth-halt.service
/lib/systemd/system/plymouth-kexec.service
/lib/systemd/system/plymouth-poweroff.service
/lib/systemd/system/plymouth-quit-wait.service
/lib/systemd/system/plymouth-quit.service
/lib/systemd/system/plymouth-read-write.service
/lib/systemd/system/plymouth-reboot.service
/lib/systemd/system/plymouth-start.service
/lib/systemd/system/plymouth-switch-root.service
/lib/systemd/system/poweroff.target.wants
/lib/systemd/system/reboot.target.wants
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/systemd-ask-password-plymouth.path
/lib/systemd/system/systemd-ask-password-plymouth.service
/sbin
/sbin/plymouthd
/usr
/usr/lib
/usr/lib/i386-linux-gnu
/usr/lib/i386-linux-gnu/plymouth
/usr/lib/i386-linux-gnu/plymouth/details.so
/usr/lib/i386-linux-gnu/plymouth/plymouth-generate-initrd
/usr/lib/i386-linux-gnu/plymouth/plymouth-populate-initrd
/usr/lib/i386-linux-gnu/plymouth/plymouth-update-initrd
/usr/lib/i386-linux-gnu/plymouth/renderers
/usr/lib/i386-linux-gnu/plymouth/renderers/drm.so
/usr/lib/i386-linux-gnu/plymouth/renderers/frame-buffer.so
/usr/lib/i386-linux-gnu/plymouth/script.so
/usr/lib/i386-linux-gnu/plymouth/text.so
/usr/lib/i386-linux-gnu/plymouth/tribar.so
/usr/sbin
/usr/sbin/plymouth-set-default-theme
/usr/share
/usr/share/doc
/usr/share/doc/plymouth
/usr/share/doc/plymouth/NEWS.Debian.gz
/usr/share/doc/plymouth/README
/usr/share/doc/plymouth/README.Debian
/usr/share/doc/plymouth/TODO
/usr/share/doc/plymouth/TODO.Debian
/usr/share/doc/plymouth/changelog.Debian.gz
/usr/share/doc/plymouth/changelog.gz
/usr/share/doc/plymouth/copyright
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/plymouth
/usr/share/initramfs-tools/scripts/panic
/usr/share/initramfs-tools/scripts/panic/plymouth
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/plymouth
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/plymouth-set-default-theme.1.gz
/usr/share/man/man1/plymouth.1.gz
/usr/share/man/man8
/usr/share/man/man8/plymouth.8.gz
/usr/share/man/man8/plymouthd.8.gz
/usr/share/plymouth
/usr/share/plymouth/debian-logo.png
/usr/share/plymouth/plymouthd.defaults
/usr/share/plymouth/themes
/usr/share/plymouth/themes/details
/usr/share/plymouth/themes/details/details.plymouth
/usr/share/plymouth/themes/text
/usr/share/plymouth/themes/text/text.plymouth
/usr/share/plymouth/themes/tribar
/usr/share/plymouth/themes/tribar/tribar.plymouth
/var
/var/lib
/var/lib/plymouth
/lib/systemd/system/halt.target.wants/plymouth-halt.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-start.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-switch-root.service
/lib/systemd/system/kexec.target.wants/plymouth-kexec.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit-wait.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit.service
/lib/systemd/system/plymouth-log.service
/lib/systemd/system/plymouth.service
/lib/systemd/system/poweroff.target.wants/plymouth-poweroff.service
/lib/systemd/system/reboot.target.wants/plymouth-reboot.service
/lib/systemd/system/sysinit.target.wants/plymouth-read-write.service
/lib/systemd/system/sysinit.target.wants/plymouth-start.service
/usr/share/man/man8/plymouth-log-viewer.8.gz
/usr/share/man/man8/plymouth-set-default-theme.8.gz
