# var/lib/dpkg/info/plymouth.conffiles
/etc/plymouth/plymouthd.conf
/etc/init.d/plymouth
