# var/lib/dpkg/info/procps.prerm
#!/bin/sh
set -e

case "$1" in
  remove|deconfigure)
	update-alternatives --remove w /usr/bin/w.procps
    ;;
  upgrade|failed-upgrade)
    ;;
esac

# Automatically added by dh_installinit
if [ -x "/etc/init.d/procps" ]; then
	invoke-rc.d procps stop || exit $?
fi
# End automatically added section


exit 0

