# var/lib/dpkg/info/pulseaudio.conffiles
/etc/pulse/default.pa
/etc/pulse/system.pa
/etc/pulse/daemon.conf
/etc/dbus-1/system.d/pulseaudio-system.conf
/etc/init.d/pulseaudio
/etc/xdg/autostart/pulseaudio.desktop
/etc/xdg/autostart/pulseaudio-kde.desktop
/etc/default/pulseaudio
