# var/lib/dpkg/info/pulseaudio.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/pulseaudio" ]; then
	invoke-rc.d pulseaudio stop || exit $?
fi
# End automatically added section
