# var/lib/dpkg/info/rsyslog.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rsyslog.service
/usr
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rsyslog
/usr/share/doc
/usr/share/doc/rsyslog
/usr/share/doc/rsyslog/changelog.Debian.gz
/usr/share/doc/rsyslog/NEWS.Debian.gz
/usr/share/doc/rsyslog/changelog.gz
/usr/share/doc/rsyslog/copyright
/usr/share/doc/rsyslog/README.Debian
/usr/share/doc/rsyslog/AUTHORS
/usr/share/man
/usr/share/man/man1
/usr/share/man/man5
/usr/share/man/man5/rsyslog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/rsyslogd.8.gz
/usr/sbin
/usr/sbin/rsyslogd
/usr/lib
/usr/lib/rsyslog
/usr/lib/rsyslog/imfile.so
/usr/lib/rsyslog/imklog.so
/usr/lib/rsyslog/immark.so
/usr/lib/rsyslog/impstats.so
/usr/lib/rsyslog/imtcp.so
/usr/lib/rsyslog/imudp.so
/usr/lib/rsyslog/imuxsock.so
/usr/lib/rsyslog/lmnet.so
/usr/lib/rsyslog/lmnetstrms.so
/usr/lib/rsyslog/lmnsd_ptcp.so
/usr/lib/rsyslog/lmregexp.so
/usr/lib/rsyslog/lmstrmsrv.so
/usr/lib/rsyslog/lmtcpclt.so
/usr/lib/rsyslog/lmtcpsrv.so
/usr/lib/rsyslog/lmzlibw.so
/usr/lib/rsyslog/mmanon.so
/usr/lib/rsyslog/mmexternal.so
/usr/lib/rsyslog/mmnormalize.so
/usr/lib/rsyslog/mmjsonparse.so
/usr/lib/rsyslog/mmutf8fix.so
/usr/lib/rsyslog/mmpstrucdata.so
/usr/lib/rsyslog/mmsequence.so
/usr/lib/rsyslog/ommail.so
/usr/lib/rsyslog/omprog.so
/usr/lib/rsyslog/omuxsock.so
/usr/lib/rsyslog/pmaixforwardedfrom.so
/usr/lib/rsyslog/pmcisconames.so
/usr/lib/rsyslog/pmlastmsg.so
/usr/lib/rsyslog/pmrfc3164sd.so
/usr/lib/rsyslog/pmsnare.so
/usr/lib/rsyslog/imkmsg.so
/usr/lib/rsyslog/imptcp.so
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/xconsole.conf
/var
/var/spool
/var/spool/rsyslog
/etc
/etc/logrotate.d
/etc/logrotate.d/rsyslog
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/rsyslog
/etc/init.d
/etc/init.d/rsyslog
/etc/default
/etc/default/rsyslog
/etc/rsyslog.conf
/etc/rsyslog.d
