# var/lib/dpkg/info/sudo.list
/.
/etc
/etc/pam.d
/etc/pam.d/sudo
/etc/sudoers.d
/etc/sudoers.d/README
/etc/init.d
/etc/init.d/sudo
/var
/var/lib
/var/lib/sudo
/usr
/usr/lib
/usr/lib/sudo
/usr/lib/sudo/sudo_noexec.so
/usr/bin
/usr/bin/sudo
/usr/bin/sudoedit
/usr/bin/sudoreplay
/usr/sbin
/usr/sbin/visudo
/usr/share
/usr/share/doc/sudo
/usr/share/doc/sudo/UPGRADE.gz
/usr/share/doc/sudo/changelog.gz
/usr/share/doc/sudo/examples
/usr/share/doc/sudo/examples/sample.sudoers
/usr/share/doc/sudo/sample.sudoers
/usr/share/doc/sudo/ChangeLog.gz
/usr/share/doc/sudo/HISTORY
/usr/share/doc/sudo/NEWS.gz
/usr/share/doc/sudo/changelog.Debian.gz
/usr/share/doc/sudo/TROUBLESHOOTING.gz
/usr/share/doc/sudo/README
/usr/share/doc/sudo/OPTIONS
/usr/share/doc/sudo/copyright
/usr/share/doc/sudo/sample.pam
/usr/share/doc/sudo/README.Debian
/usr/share/doc/sudo/NEWS.Debian.gz
/usr/share/doc/sudo/sample.syslog.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/sudo
/usr/share/man/man8
/usr/share/man/man8/sudoedit.8.gz
/usr/share/man/man8/sudoreplay.8.gz
/usr/share/man/man8/sudo.8.gz
/usr/share/man/man8/visudo.8.gz
/usr/share/man/man5
/usr/share/man/man5/sudoers.5.gz
