# var/lib/dpkg/info/systemd.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.locale1.conf
/etc/dbus-1/system.d/org.freedesktop.hostname1.conf
/etc/dbus-1/system.d/org.freedesktop.systemd1.conf
/etc/dbus-1/system.d/org.freedesktop.timedate1.conf
/etc/dbus-1/system.d/org.freedesktop.login1.conf
/etc/binfmt.d
/etc/tmpfiles.d
/etc/systemd
/etc/systemd/systemd-logind.conf
/etc/systemd/user
/etc/systemd/system.conf
/etc/systemd/user.conf
/etc/systemd/systemd-journald.conf
/etc/systemd/system
/etc/systemd/system/multi-user.target.wants
/etc/systemd/system/sysinit.target.wants
/etc/systemd/system/local-fs.target.wants
/etc/systemd/system/getty.target.wants
/etc/modules-load.d
/etc/sysctl.d
/etc/bash_completion.d
/etc/bash_completion.d/systemd-bash-completion.sh
/etc/xdg
/etc/xdg/systemd
/usr
/usr/bin
/usr/bin/systemd-cat
/usr/bin/systemd-cgtop
/usr/bin/systemd-cgls
/usr/bin/systemd-nspawn
/usr/bin/systemd-stdio-bridge
/usr/bin/systemd-analyze
/usr/lib
/usr/lib/binfmt.d
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/legacy.conf
/usr/lib/tmpfiles.d/x11.conf
/usr/lib/tmpfiles.d/debian.conf
/usr/lib/tmpfiles.d/systemd.conf
/usr/lib/tmpfiles.d/tmp.conf
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/default.target
/usr/lib/systemd/user/exit.target
/usr/lib/systemd/user/exit.service
/usr/lib/i386-linux-gnu
/usr/lib/i386-linux-gnu/systemd
/usr/lib/i386-linux-gnu/systemd/user-generators
/usr/lib/modules-load.d
/usr/lib/sysctl.d
/usr/share
/usr/share/pkgconfig
/usr/share/pkgconfig/systemd.pc
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Service.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Unit.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Automount.xml
/usr/share/dbus-1/interfaces/org.freedesktop.hostname1.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Swap.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Job.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Device.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Socket.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Path.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Mount.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Manager.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Snapshot.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Timer.xml
/usr/share/dbus-1/interfaces/org.freedesktop.timedate1.xml
/usr/share/dbus-1/interfaces/org.freedesktop.locale1.xml
/usr/share/dbus-1/interfaces/org.freedesktop.systemd1.Target.xml
/usr/share/dbus-1/services
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.locale1.service
/usr/share/dbus-1/system-services/org.freedesktop.systemd1.service
/usr/share/dbus-1/system-services/org.freedesktop.hostname1.service
/usr/share/dbus-1/system-services/org.freedesktop.timedate1.service
/usr/share/dbus-1/system-services/org.freedesktop.login1.service
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/systemd-tmpfiles.8.gz
/usr/share/man/man5
/usr/share/man/man5/systemd.socket.5.gz
/usr/share/man/man5/tmpfiles.d.5.gz
/usr/share/man/man5/systemd.target.5.gz
/usr/share/man/man5/systemd.swap.5.gz
/usr/share/man/man5/os-release.5.gz
/usr/share/man/man5/timezone.5.gz
/usr/share/man/man5/systemd.mount.5.gz
/usr/share/man/man5/binfmt.d.5.gz
/usr/share/man/man5/modules-load.d.5.gz
/usr/share/man/man5/vconsole.conf.5.gz
/usr/share/man/man5/hostname.5.gz
/usr/share/man/man5/systemd-journald.conf.5.gz
/usr/share/man/man5/systemd.conf.5.gz
/usr/share/man/man5/locale.conf.5.gz
/usr/share/man/man5/systemd-logind.conf.5.gz
/usr/share/man/man5/systemd.unit.5.gz
/usr/share/man/man5/systemd.exec.5.gz
/usr/share/man/man5/systemd.automount.5.gz
/usr/share/man/man5/systemd.path.5.gz
/usr/share/man/man5/systemd.service.5.gz
/usr/share/man/man5/sysctl.d.5.gz
/usr/share/man/man5/systemd.snapshot.5.gz
/usr/share/man/man5/systemd.device.5.gz
/usr/share/man/man5/machine-info.5.gz
/usr/share/man/man5/systemd.timer.5.gz
/usr/share/man/man5/machine-id.5.gz
/usr/share/man/man1
/usr/share/man/man1/systemd-journalctl.1.gz
/usr/share/man/man1/systemd-machine-id-setup.1.gz
/usr/share/man/man1/systemd-nspawn.1.gz
/usr/share/man/man1/systemd-cgtop.1.gz
/usr/share/man/man1/systemd-cat.1.gz
/usr/share/man/man1/systemd-notify.1.gz
/usr/share/man/man1/systemd.1.gz
/usr/share/man/man1/systemd-cgls.1.gz
/usr/share/man/man1/systemd-loginctl.1.gz
/usr/share/man/man1/systemctl.1.gz
/usr/share/man/man1/systemd-ask-password.1.gz
/usr/share/man/man3
/usr/share/man/man3/sd_seat_get_active.3.gz
/usr/share/man/man3/sd_session_is_active.3.gz
/usr/share/man/man3/sd_listen_fds.3.gz
/usr/share/man/man3/sd_login_monitor_new.3.gz
/usr/share/man/man3/sd_uid_get_state.3.gz
/usr/share/man/man3/sd_notify.3.gz
/usr/share/man/man3/sd_booted.3.gz
/usr/share/man/man3/sd_pid_get_session.3.gz
/usr/share/man/man3/sd_get_seats.3.gz
/usr/share/man/man3/sd_is_fifo.3.gz
/usr/share/man/man3/sd_readahead.3.gz
/usr/share/man/man7
/usr/share/man/man7/sd-readahead.7.gz
/usr/share/man/man7/sd-daemon.7.gz
/usr/share/man/man7/daemon.7.gz
/usr/share/man/man7/systemd.special.7.gz
/usr/share/man/man7/sd-login.7.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.login1.policy
/usr/share/polkit-1/actions/org.freedesktop.locale1.policy
/usr/share/polkit-1/actions/org.freedesktop.timedate1.policy
/usr/share/polkit-1/actions/org.freedesktop.hostname1.policy
/usr/share/polkit-1/actions/org.freedesktop.systemd1.policy
/usr/share/systemd
/usr/share/systemd/kbd-model-map
/usr/share/doc
/usr/share/doc/systemd
/usr/share/doc/systemd/changelog.Debian.gz
/usr/share/doc/systemd/copyright
/usr/share/doc/systemd/README
/usr/share/doc/systemd/README.Debian
/usr/share/doc/systemd/NEWS.gz
/usr/share/doc/systemd/DISTRO_PORTING
/bin
/bin/systemd-tty-ask-password-agent
/bin/systemd-loginctl
/bin/systemd-notify
/bin/systemd-machine-id-setup
/bin/systemctl
/bin/systemd-ask-password
/bin/systemd-journalctl
/bin/systemd-tmpfiles
/lib
/lib/lsb
/lib/lsb/init-functions.d
/lib/lsb/init-functions.d/40-systemd
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/71-seat.rules
/lib/udev/rules.d/73-seat-late.rules
/lib/udev/rules.d/70-uaccess.rules
/lib/udev/rules.d/99-systemd.rules
/lib/systemd
/lib/systemd/systemd-sysctl
/lib/systemd/systemd-remount-api-vfs
/lib/systemd/system-generators
/lib/systemd/system-generators/systemd-getty-generator
/lib/systemd/system-generators/systemd-cryptsetup-generator
/lib/systemd/systemd-multi-seat-x
/lib/systemd/systemd-readahead-collect
/lib/systemd/systemd-ac-power
/lib/systemd/systemd-cryptsetup
/lib/systemd/systemd-modules-load
/lib/systemd/systemd-uaccess
/lib/systemd/systemd-hostnamed
/lib/systemd/systemd-fsck
/lib/systemd/systemd-readahead-replay
/lib/systemd/systemd-timestamp
/lib/systemd/systemd-journald
/lib/systemd/systemd-detect-virt
/lib/systemd/system-shutdown
/lib/systemd/systemd
/lib/systemd/systemd-quotacheck
/lib/systemd/systemd-timedated
/lib/systemd/systemd-initctl
/lib/systemd/systemd-shutdown
/lib/systemd/systemd-cgroups-agent
/lib/systemd/systemd-random-seed
/lib/systemd/systemd-user-sessions
/lib/systemd/debian-fixup
/lib/systemd/systemd-shutdownd
/lib/systemd/systemd-update-utmp
/lib/systemd/systemd-logind
/lib/systemd/systemd-vconsole-setup
/lib/systemd/systemd-reply-password
/lib/systemd/systemd-localed
/lib/systemd/system
/lib/systemd/system/time-sync.target
/lib/systemd/system/run-user.mount
/lib/systemd/system/umount.target
/lib/systemd/system/poweroff.target
/lib/systemd/system/user@.service
/lib/systemd/system/systemd-journald.service
/lib/systemd/system/systemd-tmpfiles-setup.service
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/proc-sys-fs-binfmt_misc.automount
/lib/systemd/system/shutdown.target.wants
/lib/systemd/system/debian-fixup.service
/lib/systemd/system/systemd-random-seed-load.service
/lib/systemd/system/systemd-modules-load.service
/lib/systemd/system/halt.service
/lib/systemd/system/sys-fs-fuse-connections.mount
/lib/systemd/system/sysinit.target
/lib/systemd/system/kexec.service
/lib/systemd/system/runlevel2.target.wants
/lib/systemd/system/systemd-ask-password-console.path
/lib/systemd/system/systemd-readahead-replay.service
/lib/systemd/system/http-daemon.target
/lib/systemd/system/systemd-remount-api-vfs.service
/lib/systemd/system/systemd-shutdownd.service
/lib/systemd/system/emergency.target
/lib/systemd/system/fsck@.service
/lib/systemd/system/rpcbind.target
/lib/systemd/system/mail-transfer-agent.target
/lib/systemd/system/printer.target
/lib/systemd/system/systemd-journald.socket
/lib/systemd/system/systemd-timedated.service
/lib/systemd/system/nss-lookup.target
/lib/systemd/system/systemd-readahead-done.timer
/lib/systemd/system/run-lock.mount
/lib/systemd/system/systemd-tmpfiles-clean.service
/lib/systemd/system/bluetooth.target
/lib/systemd/system/rescue.service
/lib/systemd/system/syslog.target
/lib/systemd/system/systemd-logind.service
/lib/systemd/system/var-lock.mount
/lib/systemd/system/swap.target
/lib/systemd/system/runlevel3.target.wants
/lib/systemd/system/dev-hugepages.mount
/lib/systemd/system/systemd-readahead-collect.service
/lib/systemd/system/systemd-binfmt.service
/lib/systemd/system/runlevel1.target.wants
/lib/systemd/system/systemd-readahead-done.service
/lib/systemd/system/remote-fs-pre.target
/lib/systemd/system/rescue.target
/lib/systemd/system/kexec.target
/lib/systemd/system/sigpwr.target
/lib/systemd/system/runlevel5.target.wants
/lib/systemd/system/poweroff.service
/lib/systemd/system/sys-kernel-config.mount
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/cryptsetup.target
/lib/systemd/system/systemd-update-utmp-shutdown.service
/lib/systemd/system/dev-mqueue.mount
/lib/systemd/system/getty@.service
/lib/systemd/system/systemd-ask-password-wall.path
/lib/systemd/system/basic.target.wants
/lib/systemd/system/sound.target
/lib/systemd/system/syslog.socket
/lib/systemd/system/local-fs.target.wants
/lib/systemd/system/var-run.mount
/lib/systemd/system/console-shell.service
/lib/systemd/system/reboot.service
/lib/systemd/system/local-fs-pre.target
/lib/systemd/system/systemd-shutdownd.socket
/lib/systemd/system/smartcard.target
/lib/systemd/system/systemd-vconsole-setup.service
/lib/systemd/system/graphical.target.wants
/lib/systemd/system/sys-kernel-security.mount
/lib/systemd/system/serial-getty@.service
/lib/systemd/system/systemd-user-sessions.service
/lib/systemd/system/remount-rootfs.service
/lib/systemd/system/proc-sys-fs-binfmt_misc.mount
/lib/systemd/system/shutdown.target
/lib/systemd/system/sockets.target
/lib/systemd/system/getty.target
/lib/systemd/system/halt.target
/lib/systemd/system/basic.target
/lib/systemd/system/systemd-initctl.socket
/lib/systemd/system/systemd-ask-password-console.service
/lib/systemd/system/systemd-hostnamed.service
/lib/systemd/system/remote-fs.target
/lib/systemd/system/emergency.service
/lib/systemd/system/final.target
/lib/systemd/system/systemd-tmpfiles-clean.timer
/lib/systemd/system/graphical.target
/lib/systemd/system/reboot.target
/lib/systemd/system/quotacheck.service
/lib/systemd/system/fsck-root.service
/lib/systemd/system/systemd-sysctl.service
/lib/systemd/system/systemd-random-seed-save.service
/lib/systemd/system/sockets.target.wants
/lib/systemd/system/multi-user.target
/lib/systemd/system/sys-kernel-debug.mount
/lib/systemd/system/systemd-update-utmp-runlevel.service
/lib/systemd/system/ifup@.service
/lib/systemd/system/network.target
/lib/systemd/system/quotaon.service
/lib/systemd/system/systemd-ask-password-wall.service
/lib/systemd/system/runlevel4.target.wants
/lib/systemd/system/local-fs.target
/lib/systemd/system/systemd-initctl.service
/lib/systemd/system/systemd-localed.service
/lib/systemd/systemd-binfmt
/var
/var/lib
/var/lib/systemd
/var/log
/var/log/journal
/etc/systemd/system/multi-user.target.wants/remote-fs.target
/etc/systemd/system/getty.target.wants/getty@tty1.service
/etc/modules-load.d/modules.conf
/etc/xdg/systemd/user
/usr/lib/systemd/user/printer.target
/usr/lib/systemd/user/bluetooth.target
/usr/lib/systemd/user/sound.target
/usr/lib/systemd/user/shutdown.target
/usr/lib/systemd/user/sockets.target
/usr/share/dbus-1/services/org.freedesktop.systemd1.service
/usr/share/man/man3/sd_uid_get_sessions.3.gz
/usr/share/man/man3/sd_is_socket_unix.3.gz
/usr/share/man/man3/sd_session_get_class.3.gz
/usr/share/man/man3/sd_session_get_display.3.gz
/usr/share/man/man3/sd_seat_can_multi_session.3.gz
/usr/share/man/man3/sd_is_mq.3.gz
/usr/share/man/man3/sd_session_get_seat.3.gz
/usr/share/man/man3/sd_uid_get_seats.3.gz
/usr/share/man/man3/sd_session_get_uid.3.gz
/usr/share/man/man3/sd_login_monitor_get_fd.3.gz
/usr/share/man/man3/sd_get_uids.3.gz
/usr/share/man/man3/sd_pid_get_owner_uid.3.gz
/usr/share/man/man3/sd_is_socket_inet.3.gz
/usr/share/man/man3/sd_login_monitor_flush.3.gz
/usr/share/man/man3/sd_login_monitor_unref.3.gz
/usr/share/man/man3/sd_uid_is_on_seat.3.gz
/usr/share/man/man3/sd_get_sessions.3.gz
/usr/share/man/man3/sd_seat_get_sessions.3.gz
/usr/share/man/man3/sd_is_socket.3.gz
/usr/share/man/man3/sd_session_get_type.3.gz
/usr/share/man/man3/sd_notifyf.3.gz
/usr/share/man/man3/sd_session_get_service.3.gz
/usr/share/man/man3/sd_pid_get_unit.3.gz
/bin/systemd
/lib/systemd/system/runlevel1.target
/lib/systemd/system/multi-user.target.wants/systemd-logind.service
/lib/systemd/system/multi-user.target.wants/systemd-ask-password-wall.path
/lib/systemd/system/multi-user.target.wants/systemd-user-sessions.service
/lib/systemd/system/multi-user.target.wants/getty.target
/lib/systemd/system/shutdown.target.wants/systemd-update-utmp-shutdown.service
/lib/systemd/system/shutdown.target.wants/systemd-random-seed-save.service
/lib/systemd/system/procps.service
/lib/systemd/system/hwclock.service
/lib/systemd/system/runlevel2.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/single.service
/lib/systemd/system/runlevel4.target
/lib/systemd/system/stop-bootlogd-single.service
/lib/systemd/system/mountall-bootclean.service
/lib/systemd/system/mountnfs.service
/lib/systemd/system/dbus-org.freedesktop.login1.service
/lib/systemd/system/dbus-org.freedesktop.hostname1.service
/lib/systemd/system/cryptdisks.service
/lib/systemd/system/default.target
/lib/systemd/system/hostname.service
/lib/systemd/system/hwclockfirst.service
/lib/systemd/system/udev-mtab.service
/lib/systemd/system/module-init-tools.service
/lib/systemd/system/checkroot-bootclean.service
/lib/systemd/system/runlevel2.target
/lib/systemd/system/checkroot.service
/lib/systemd/system/autovt@.service
/lib/systemd/system/rmnologin.service
/lib/systemd/system/checkfs.service
/lib/systemd/system/runlevel3.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/x11-common.service
/lib/systemd/system/mountall.service
/lib/systemd/system/fuse.service
/lib/systemd/system/runlevel1.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/mountkernfs.service
/lib/systemd/system/runlevel5.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/runlevel0.target
/lib/systemd/system/bootmisc.service
/lib/systemd/system/sysinit.target.wants/systemd-journald.service
/lib/systemd/system/sysinit.target.wants/systemd-tmpfiles-setup.service
/lib/systemd/system/sysinit.target.wants/proc-sys-fs-binfmt_misc.automount
/lib/systemd/system/sysinit.target.wants/debian-fixup.service
/lib/systemd/system/sysinit.target.wants/systemd-random-seed-load.service
/lib/systemd/system/sysinit.target.wants/systemd-modules-load.service
/lib/systemd/system/sysinit.target.wants/sys-fs-fuse-connections.mount
/lib/systemd/system/sysinit.target.wants/systemd-ask-password-console.path
/lib/systemd/system/sysinit.target.wants/dev-hugepages.mount
/lib/systemd/system/sysinit.target.wants/systemd-binfmt.service
/lib/systemd/system/sysinit.target.wants/sys-kernel-config.mount
/lib/systemd/system/sysinit.target.wants/cryptsetup.target
/lib/systemd/system/sysinit.target.wants/dev-mqueue.mount
/lib/systemd/system/sysinit.target.wants/sys-kernel-security.mount
/lib/systemd/system/sysinit.target.wants/systemd-sysctl.service
/lib/systemd/system/sysinit.target.wants/sys-kernel-debug.mount
/lib/systemd/system/bootlogd.service
/lib/systemd/system/runlevel6.target
/lib/systemd/system/basic.target.wants/systemd-tmpfiles-clean.timer
/lib/systemd/system/killprocs.service
/lib/systemd/system/mountdevsubfs.service
/lib/systemd/system/local-fs.target.wants/run-user.mount
/lib/systemd/system/local-fs.target.wants/systemd-remount-api-vfs.service
/lib/systemd/system/local-fs.target.wants/run-lock.mount
/lib/systemd/system/local-fs.target.wants/var-lock.mount
/lib/systemd/system/local-fs.target.wants/var-run.mount
/lib/systemd/system/local-fs.target.wants/remount-rootfs.service
/lib/systemd/system/local-fs.target.wants/fsck-root.service
/lib/systemd/system/mtab.service
/lib/systemd/system/stop-bootlogd.service
/lib/systemd/system/kmod.service
/lib/systemd/system/mountnfs-bootclean.service
/lib/systemd/system/ctrl-alt-del.target
/lib/systemd/system/cryptdisks-early.service
/lib/systemd/system/dbus-org.freedesktop.locale1.service
/lib/systemd/system/sockets.target.wants/systemd-journald.socket
/lib/systemd/system/sockets.target.wants/systemd-shutdownd.socket
/lib/systemd/system/sockets.target.wants/systemd-initctl.socket
/lib/systemd/system/dbus-org.freedesktop.timedate1.service
/lib/systemd/system/runlevel3.target
/lib/systemd/system/runlevel4.target.wants/systemd-update-utmp-runlevel.service
/lib/systemd/system/runlevel5.target
/lib/systemd/system/urandom.service
