# var/lib/dpkg/info/udev.list
/.
/sbin
/usr
/usr/share
/usr/share/bug
/usr/share/bug/udev
/usr/share/doc
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/udev
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/udev
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/udev
/usr/share/pkgconfig
/usr/share/pkgconfig/udev.pc
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/udevadm.8.gz
/usr/share/man/man8/systemd-udevd.service.8.gz
/usr/share/man/man7
/usr/share/man/man7/udev.7.gz
/bin
/bin/udevadm
/lib
/lib/systemd
/lib/systemd/systemd-udevd
/lib/systemd/system
/lib/systemd/system/udev-finish.service
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/sockets.target.wants
/lib/systemd/system/systemd-udevd.service
/lib/systemd/system/systemd-udevd-kernel.socket
/lib/systemd/system/systemd-udevd-control.socket
/lib/systemd/system/systemd-udev-trigger.service
/lib/systemd/system/systemd-udev-settle.service
/lib/udev
/lib/udev/accelerometer
/lib/udev/ata_id
/lib/udev/cdrom_id
/lib/udev/collect
/lib/udev/mtd_probe
/lib/udev/scsi_id
/lib/udev/v4l_id
/lib/udev/write_net_rules
/lib/udev/udev-finish
/lib/udev/net.agent
/lib/udev/logger.agent
/lib/udev/dsl-modem.agent
/lib/udev/rule_generator.functions
/lib/udev/hotplug.functions
/lib/udev/rules.d
/lib/udev/rules.d/80-networking.rules
/lib/udev/rules.d/75-persistent-net-generator.rules
/lib/udev/rules.d/73-idrac.rules
/lib/udev/rules.d/70-power-switch.rules
/lib/udev/rules.d/75-probe_mtd.rules
/lib/udev/rules.d/61-accelerometer.rules
/lib/udev/rules.d/60-persistent-v4l.rules
/lib/udev/rules.d/60-cdrom_id.rules
/lib/udev/rules.d/80-drivers.rules
/lib/udev/rules.d/50-firmware.rules
/lib/udev/rules.d/95-udev-late.rules
/lib/udev/rules.d/80-net-setup-link.rules
/lib/udev/rules.d/78-sound-card.rules
/lib/udev/rules.d/75-tty-description.rules
/lib/udev/rules.d/75-net-description.rules
/lib/udev/rules.d/64-btrfs.rules
/lib/udev/rules.d/60-persistent-storage.rules
/lib/udev/rules.d/60-persistent-alsa.rules
/lib/udev/rules.d/60-persistent-input.rules
/lib/udev/rules.d/60-persistent-serial.rules
/lib/udev/rules.d/60-persistent-storage-tape.rules
/lib/udev/rules.d/60-keyboard.rules
/lib/udev/rules.d/60-drm.rules
/lib/udev/rules.d/50-udev-default.rules
/lib/udev/rules.d/42-usb-hid-pm.rules
/lib/udev/hwdb.d
/lib/udev/hwdb.d/60-keyboard.hwdb
/lib/udev/hwdb.d/20-net-ifname.hwdb
/lib/udev/hwdb.d/20-OUI.hwdb
/lib/udev/hwdb.d/20-acpi-vendor.hwdb
/lib/udev/hwdb.d/20-bluetooth-vendor-product.hwdb
/lib/udev/hwdb.d/20-sdio-classes.hwdb
/lib/udev/hwdb.d/20-sdio-vendor-model.hwdb
/lib/udev/hwdb.d/20-usb-classes.hwdb
/lib/udev/hwdb.d/20-usb-vendor-model.hwdb
/lib/udev/hwdb.d/20-pci-classes.hwdb
/lib/udev/hwdb.d/20-pci-vendor-model.hwdb
/etc
/etc/init.d
/etc/init.d/udev-finish
/etc/init.d/udev
/etc/init
/etc/init/udev-fallback-graphics.conf
/etc/init/udevmonitor.conf
/etc/init/udevtrigger.conf
/etc/init/udev-finish.conf
/etc/init/udev.conf
/etc/modprobe.d
/etc/modprobe.d/fbdev-blacklist.conf
/etc/udev
/etc/udev/hwdb.d
/etc/udev/rules.d
/etc/udev/udev.conf
/sbin/udevd
/sbin/udevadm
/usr/share/doc/udev
/usr/share/man/man8/udevd.8.gz
/usr/share/man/man8/systemd-udevd.8.gz
/usr/share/man/man8/systemd-udevd-kernel.socket.8.gz
/usr/share/man/man8/systemd-udevd-control.socket.8.gz
/lib/systemd/system/udev.service
/lib/systemd/system/sysinit.target.wants/udev-finish.service
/lib/systemd/system/sysinit.target.wants/systemd-udevd.service
/lib/systemd/system/sysinit.target.wants/systemd-udev-trigger.service
/lib/systemd/system/sockets.target.wants/systemd-udevd-kernel.socket
/lib/systemd/system/sockets.target.wants/systemd-udevd-control.socket
