# var/lib/dpkg/info/udev.conffiles
/etc/init.d/udev-mtab
/etc/init.d/udev
/etc/udev/udev.conf
/etc/udev/links.conf
/etc/modprobe.d/blacklist.conf
/etc/modprobe.d/fbdev-blacklist.conf
