# var/lib/dpkg/info/udev.conffiles
/etc/init.d/udev
/etc/init/udev.conf
/etc/init/udevmonitor.conf
/etc/init/udevtrigger.conf
/etc/udev/udev.conf
