# var/lib/dpkg/info/udev.list
/.
/usr
/usr/lib
/usr/lib/ConsoleKit
/usr/lib/ConsoleKit/run-seat.d
/usr/share
/usr/share/bug
/usr/share/bug/udev
/usr/share/pkgconfig
/usr/share/pkgconfig/udev.pc
/usr/share/doc
/usr/share/doc/udev
/usr/share/doc/udev/NEWS.Debian.gz
/usr/share/doc/udev/changelog.gz
/usr/share/doc/udev/examples
/usr/share/doc/udev/examples/udev.vim
/usr/share/doc/udev/NEWS.gz
/usr/share/doc/udev/changelog.Debian.gz
/usr/share/doc/udev/copyright
/usr/share/doc/udev/TODO
/usr/share/doc/udev/README.keymap.txt.gz
/usr/share/doc/udev/README.Debian.gz
/usr/share/doc/udev/README.gz
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/udev
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/udev
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/udev
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/udev.7.gz
/usr/share/man/man8
/usr/share/man/man8/udevadm.8.gz
/usr/share/man/man8/udevd.8.gz
/usr/share/man/man8/scsi_id.8.gz
/lib
/lib/firmware
/lib/udev
/lib/udev/logger.agent
/lib/udev/v4l_id
/lib/udev/keymaps
/lib/udev/keymaps/hewlett-packard-tx2
/lib/udev/keymaps/module-asus-w3j
/lib/udev/keymaps/toshiba-satellite_a100
/lib/udev/keymaps/lenovo-thinkpad_x6_tablet
/lib/udev/keymaps/acer-aspire_5920g
/lib/udev/keymaps/samsung-sq1us
/lib/udev/keymaps/acer-aspire_6920
/lib/udev/keymaps/medionnb-a555
/lib/udev/keymaps/olpc-xo
/lib/udev/keymaps/hewlett-packard
/lib/udev/keymaps/module-ibm
/lib/udev/keymaps/lenovo-ideapad
/lib/udev/keymaps/fujitsu-amilo_pa_2548
/lib/udev/keymaps/hewlett-packard-compaq_elitebook
/lib/udev/keymaps/zepto-znote
/lib/udev/keymaps/onkyo
/lib/udev/keymaps/acer-aspire_5720
/lib/udev/keymaps/hewlett-packard-presario-2100
/lib/udev/keymaps/logitech-wave-pro-cordless
/lib/udev/keymaps/lenovo-3000
/lib/udev/keymaps/fujitsu-esprimo_mobile_v6
/lib/udev/keymaps/fujitsu-amilo_si_1520
/lib/udev/keymaps/module-sony
/lib/udev/keymaps/lenovo-thinkpad_x200_tablet
/lib/udev/keymaps/inventec-symphony_6.0_7.0
/lib/udev/keymaps/logitech-wave
/lib/udev/keymaps/hewlett-packard-2510p_2530p
/lib/udev/keymaps/force-release
/lib/udev/keymaps/force-release/hp-other
/lib/udev/keymaps/force-release/common-volume-keys
/lib/udev/keymaps/force-release/dell-touchpad
/lib/udev/keymaps/force-release/samsung-other
/lib/udev/keymaps/compaq-e_evo
/lib/udev/keymaps/hewlett-packard-tablet
/lib/udev/keymaps/maxdata-pro_7000
/lib/udev/keymaps/toshiba-satellite_a110
/lib/udev/keymaps/asus
/lib/udev/keymaps/samsung-sx20s
/lib/udev/keymaps/lg-x110
/lib/udev/keymaps/module-sony-vgn
/lib/udev/keymaps/acer
/lib/udev/keymaps/hewlett-packard-pavilion
/lib/udev/keymaps/fujitsu-amilo_pro_v3205
/lib/udev/keymaps/samsung-other
/lib/udev/keymaps/acer-travelmate_c300
/lib/udev/keymaps/everex-xt5000
/lib/udev/keymaps/micro-star
/lib/udev/keymaps/ibm-thinkpad-usb-keyboard-trackpoint
/lib/udev/keymaps/fujitsu-amilo_pro_edition_v3505
/lib/udev/keymaps/module-lenovo
/lib/udev/keymaps/module-sony-old
/lib/udev/keymaps/toshiba-satellite_m30x
/lib/udev/keymaps/oqo-model2
/lib/udev/keymaps/fujitsu-esprimo_mobile_v5
/lib/udev/keymaps/logitech-wave-cordless
/lib/udev/keymaps/medion-fid2060
/lib/udev/keymaps/lenovo-thinkpad-usb-keyboard-trackpoint
/lib/udev/keymaps/dell
/lib/udev/dsl-modem.agent
/lib/udev/write_cd_rules
/lib/udev/input_id
/lib/udev/hotplug.functions
/lib/udev/rule_generator.functions
/lib/udev/path_id
/lib/udev/keyboard-force-release.sh
/lib/udev/create_static_nodes
/lib/udev/findkeyboards
/lib/udev/hid2hci
/lib/udev/scsi_id
/lib/udev/edd_id
/lib/udev/pci-db
/lib/udev/mobile-action-modeswitch
/lib/udev/devices
/lib/udev/net.agent
/lib/udev/write_net_rules
/lib/udev/udev-acl
/lib/udev/ata_id
/lib/udev/mtd_probe
/lib/udev/firmware.agent
/lib/udev/collect
/lib/udev/usb-db
/lib/udev/rules.d
/lib/udev/rules.d/50-udev-default.rules
/lib/udev/rules.d/70-acl.rules
/lib/udev/rules.d/60-persistent-input.rules
/lib/udev/rules.d/61-mobile-action.rules
/lib/udev/rules.d/78-sound-card.rules
/lib/udev/rules.d/75-probe_mtd.rules
/lib/udev/rules.d/60-persistent-v4l.rules
/lib/udev/rules.d/91-permissions.rules
/lib/udev/rules.d/60-persistent-alsa.rules
/lib/udev/rules.d/75-net-description.rules
/lib/udev/rules.d/60-persistent-storage.rules
/lib/udev/rules.d/70-hid2hci.rules
/lib/udev/rules.d/75-persistent-net-generator.rules
/lib/udev/rules.d/80-drivers.rules
/lib/udev/rules.d/60-persistent-storage-tape.rules
/lib/udev/rules.d/95-keyboard-force-release.rules
/lib/udev/rules.d/79-fstab_import.rules
/lib/udev/rules.d/95-udev-late.rules
/lib/udev/rules.d/75-cd-aliases-generator.rules
/lib/udev/rules.d/75-tty-description.rules
/lib/udev/rules.d/60-persistent-serial.rules
/lib/udev/rules.d/95-keymap.rules
/lib/udev/usb_id
/lib/udev/keymap
/lib/udev/cdrom_id
/lib/udev/fstab_import
/sbin
/sbin/udevadm
/sbin/udevd
/etc
/etc/init.d
/etc/init.d/udev-mtab
/etc/init.d/udev
/etc/udev
/etc/udev/udev.conf
/etc/udev/rules.d
/etc/udev/.dev
/etc/udev/links.conf
/etc/modprobe.d
/etc/modprobe.d/blacklist.conf
/etc/modprobe.d/fbdev-blacklist.conf
/usr/lib/ConsoleKit/run-seat.d/udev-acl.ck
