# var/lib/dpkg/info/util-linux.conffiles
/etc/init.d/hwclock.sh
/etc/init.d/hwclockfirst.sh
