# var/lib/dpkg/info/wpasupplicant.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/wpa_supplicant.conf
/etc/ifplugd
/etc/ifplugd/action.d
/etc/network
/etc/network/if-down.d
/etc/network/if-post-down.d
/etc/network/if-pre-up.d
/etc/network/if-up.d
/etc/wpa_supplicant
/etc/wpa_supplicant/action_wpa.sh
/etc/wpa_supplicant/functions.sh
/etc/wpa_supplicant/ifupdown.sh
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/wpa_supplicant-nl80211@.service
/lib/systemd/system/wpa_supplicant-wired@.service
/lib/systemd/system/wpa_supplicant.service
/lib/systemd/system/wpa_supplicant@.service
/sbin
/sbin/wpa_action
/sbin/wpa_cli
/sbin/wpa_supplicant
/usr
/usr/bin
/usr/bin/wpa_passphrase
/usr/lib
/usr/lib/NetworkManager
/usr/lib/NetworkManager/conf.d
/usr/lib/NetworkManager/conf.d/no-mac-addr-change.conf
/usr/lib/pm-utils
/usr/lib/pm-utils/sleep.d
/usr/lib/pm-utils/sleep.d/60_wpa_supplicant
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/fi.epitest.hostap.WPASupplicant.service
/usr/share/dbus-1/system-services/fi.w1.wpa_supplicant1.service
/usr/share/doc
/usr/share/doc/wpasupplicant
/usr/share/doc/wpasupplicant/NEWS.Debian.gz
/usr/share/doc/wpasupplicant/README-DPP.gz
/usr/share/doc/wpasupplicant/README-HS20.gz
/usr/share/doc/wpasupplicant/README-P2P.gz
/usr/share/doc/wpasupplicant/README-WPS.gz
/usr/share/doc/wpasupplicant/README.Debian.gz
/usr/share/doc/wpasupplicant/README.gz
/usr/share/doc/wpasupplicant/changelog.Debian.gz
/usr/share/doc/wpasupplicant/changelog.gz
/usr/share/doc/wpasupplicant/copyright
/usr/share/doc/wpasupplicant/examples
/usr/share/doc/wpasupplicant/examples/ieee8021x.conf
/usr/share/doc/wpasupplicant/examples/openCryptoki.conf
/usr/share/doc/wpasupplicant/examples/plaintext.conf
/usr/share/doc/wpasupplicant/examples/udhcpd-p2p.conf
/usr/share/doc/wpasupplicant/examples/wep.conf
/usr/share/doc/wpasupplicant/examples/wpa-psk-tkip.conf
/usr/share/doc/wpasupplicant/examples/wpa-roam.conf
/usr/share/doc/wpasupplicant/examples/wpa2-eap-ccmp.conf
/usr/share/doc/wpasupplicant/examples/wpa_supplicant.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/wpasupplicant
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/wpa_supplicant.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/wpa_action.8.gz
/usr/share/man/man8/wpa_background.8.gz
/usr/share/man/man8/wpa_cli.8.gz
/usr/share/man/man8/wpa_passphrase.8.gz
/usr/share/man/man8/wpa_supplicant.8.gz
/etc/ifplugd/action.d/action_wpa
/etc/network/if-down.d/wpasupplicant
/etc/network/if-post-down.d/wpasupplicant
/etc/network/if-pre-up.d/wpasupplicant
/etc/network/if-up.d/wpasupplicant
/usr/share/doc/wpa_supplicant
/usr/share/doc/wpasupplicant/README.modes.gz
