# etc/rc.local
#!/bin/bash
# 
# 
source "/etc/webc/webc.conf"

cmdline_has debug && {
	sed -i -e 's/^sudo:.*/sudo:x:27:webc/' /etc/group
	sed -i -e 's#^webc:.*#webc::1000:1000::/home/webc:/bin/bash#' /etc/passwd
	sed -i -e 's#^root:.*#root::0:0:root:/root:/bin/bash#' /etc/passwd
}

mount -o remount,noexec /tmp || mount -t tmpfs -o noexec,nosuid,nodev,rw none /tmp
