# etc/systemd/system/webc.service
[Unit]
Description=Service to setup Webconverger as a superuser

[Service]
ExecStart=/etc/webc/live-config.sh
Restart=always
