# etc/network/if-up.d/ping
#!/bin/sh
. /etc/webc/webc.conf
server="http://ping.webconverger.org/"
sleeping=5

webcping () {
	logger "Pinging $server"
	wget -qO- --timeout=5 --post-data="M=$( 
		mac_address | md5sum | awk '{print $1}'
	)&V=${webc_version}" "$server"
}

runonce=/var/run/webcping
mkdir /var/lock/webcping 2>/dev/null || exit 0

while ! test -f "$runonce"
do
	if webcping
	then
		logger "Successful ping to $server"
		date > "$runonce"
		exit 0
	fi

	logger "Failed to ping $server, sleeping $sleeping seconds"
	sleep "$sleeping"

done &
