# usr/share/insserv/overrides/devfsd
### BEGIN INIT INFO
# Provides:          devfsd
# Required-Start:    mountdevsubfs 
# Required-Stop:     reboot
# Default-Start:     S 1 2 3 4 5
# Default-Stop:      0 6
### END INIT INFO
# Reported to BTS as bug #324671
