# usr/share/insserv/overrides/hotplug-net
### BEGIN INIT INFO
# Provides:          hotplug-net
# Required-Start:    $local_fs ifupdown
# Required-Stop:     $local_fs
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
