# usr/share/insserv/overrides/hotplug
### BEGIN INIT INFO
# Provides:          hotplug
# Required-Start:    mountdevsubfs checkroot $local_fs
# Required-Stop:     $local_fs
# Default-Start:     S
# Default-Stop:      0 6
### END INIT INFO
