# usr/share/insserv/overrides/modutils
### BEGIN INIT INFO
# Provides:          modutils
# Required-Start:    checkroot
# Required-Stop:
# Default-Start:     S
# Default-Stop:
### END INIT INFO
# Reported to BTS as #335311
